# acme_test_server/acme.service
[Unit]
Description=ACME Test Server
Requires=acme.socket

[Service]
Type=exec
ExecStart=/usr/bin/acme_test_server
Restart=on-failure
